# roles/IMAP/files/etc/cron.d/doveadm
MAILTO=root

7 5 * * *  vmail  test -x /usr/bin/doveadm && nice -n 19 /usr/bin/doveadm sis deduplicate /home/mail/attachments /home/mail/attachments/queue
