# roles/IMAP/files/etc/systemd/system/dovecot-auth-proxy.socket
[Socket]
SocketUser=dovecot
SocketGroup=dovecot
SocketMode=0600
ListenStream=/run/dovecot/auth-proxy

[Install]
WantedBy=sockets.target
