# roles/common/files/lib/systemd/system/bacula-fd.service
[Unit]
Description=Bacula File Daemon service
After=network.target

[Service]
Type=forking
PIDFile=/var/run/bacula/bacula-fd.9112.pid
StandardOutput=syslog
ExecStart=/usr/sbin/bacula-fd -c /etc/bacula/bacula-fd.conf

# Hardening
NoNewPrivileges=yes
PrivateDevices=yes
ProtectHome=read-only
ProtectSystem=full
PrivateTmp=yes
ReadOnlyDirectories=/
ReadWriteDirectories=-/var/lib
ReadWriteDirectories=-/var/run/bacula

[Install]
WantedBy=multi-user.target
