# roles/nextcloud/files/etc/cron.d/nextcloud
MAILTO=root
*/5 * * * * _nextcloud php -f /usr/local/share/nextcloud/cron.php
