# roles/common/files/etc/systemd/system/bacula-fd.service
[Unit]
Description=Bacula File Daemon service
After=network.target

[Service]
Type=simple
StandardOutput=syslog
ExecStart=/usr/sbin/bacula-fd -f -c /etc/bacula/bacula-fd.conf

# Hardening
NoNewPrivileges=yes
PrivateDevices=yes
ProtectHome=read-only
ProtectSystem=strict
PrivateTmp=yes
ReadWriteDirectories=-/var/lib
ReadWriteDirectories=-/var/run/bacula

[Install]
WantedBy=multi-user.target
